# deploy/n-audit-sentinel.service
[Unit]
Description=N-Audit Sentinel (systemd service for local audits)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/n-audit
Environment=LOG_FILE=/var/lib/n-audit/session.log
ExecStart=/usr/local/bin/n-audit-sentinel --system-audit --log-file ${LOG_FILE} --sign-key /var/lib/n-audit/signing/id_ed25519
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

# Notes:
# - Create the signing directory and key: mkdir -p /var/lib/n-audit/signing && ssh-keygen -t ed25519 -f /var/lib/n-audit/signing/id_ed25519
# - Protect the key: chmod 700 /var/lib/n-audit/signing && chmod 600 /var/lib/n-audit/signing/id_ed25519
# - Consider adjusting `User=` to a non-root user and grant required capabilities if desired.