# deploy/n-audit-sentinel.service
[Unit]
Description=N-Audit Sentinel (systemd service for local audits)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=n-audit
WorkingDirectory=/var/lib/n-audit
Environment=LOG_FILE=/var/lib/n-audit/session.log
ExecStart=/usr/local/bin/n-audit-sentinel --system-audit --log-file ${LOG_FILE} --sign-key /var/lib/n-audit/signing/id_ed25519
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
ProtectSystem=full
ProtectHome=true
PrivateTmp=yes
NoNewPrivileges=yes
ReadOnlyPaths=/usr
RestrictAddressFamilies=AF_INET AF_INET6
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target

# Notes:
# - Create the `n-audit` service user and ensure directories are owned by it:
#   sudo useradd --system --no-create-home --shell /sbin/nologin n-audit || true
#   sudo mkdir -p /var/lib/n-audit/signing && sudo chown -R n-audit:n-audit /var/lib/n-audit
# - Generate signing key: sudo -u n-audit ssh-keygen -t ed25519 -N "" -f /var/lib/n-audit/signing/id_ed25519
# - Protect the key: chmod 700 /var/lib/n-audit/signing && chmod 600 /var/lib/n-audit/signing/id_ed25519
# - If additional capabilities are required, grant only the minimal set via `AmbientCapabilities` or adjust the unit carefully.